# systemd/system/flarum-log-cleanup.timer
[Unit]
Description=Run flarum log cleanup every day

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=timers.target
