# systemd/system/stopforumspam-update-hourly.timer
[Unit]
Description=Stop Forum Spam hourly data update

[Timer]
OnCalendar=hourly
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
