# systemd/system/stopforumspam-update-daily.timer
[Unit]
Description=Stop Forum Spam daily data update

[Timer]
OnCalendar=daily
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
