# systemd/system/mariadb.service.d/local.conf
[Service]
Restart=always
RestartSec=100ms
SupplementaryGroups=io_uring
