# systemd/system/geoipupdate.service.d/local.conf
[Service]
CPUSchedulingPolicy=batch
IOSchedulingPriority=7
Nice=19
PrivateIPC=true
ProcSubset=pid
User=geoipupdate
